# Management_system/Source_code/services_systemd/ControlPredictivo.timer
[Unit]
Description=run my script

[Timer]
OnCalendar=*-*-* 00/2:3:00
Persistent=true

[Install]
WantedBy=timers.target
